# exchanger/cron/cron.d
*/10 * * * * /var/www/html/cron/command.sh
